# server/cron/crontab
*/5 *  *   *   1-5     /home/deton/bin/syuttaikin.py >/dev/null
